# xremap.service.template
[Unit]
Description=xremap

[Service]
ExecStart=/path/to/xremap /path/to/config.yml
Restart=always

[Install]
WantedBy=default.target
